# monitoring/ldap-repl-check.service
[Unit]
Description=Chequeo de retraso de replicacion LDAP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ldap_repl_check.sh
